# mpv.service
[Unit]
Description=MPV background service
[Service]
ExecStart=/usr/bin/mpv -idle=once
Type=simple
[Install]
WantedBy=default.target
